# roles/paperclip/templates/paperclip.service.j2
[Unit]
Description=Paperclip AI Agent Orchestration
After=network.target

[Service]
Type=simple
User={{ paperclip_user }}
Group={{ paperclip_user }}
WorkingDirectory={{ paperclip_home }}

# Environment
Environment="NODE_ENV=production"
Environment="PORT={{ paperclip_port }}"
Environment="HOME={{ paperclip_home }}"
Environment="PATH={{ npm_prefix.stdout | trim }}/bin:/usr/local/bin:/usr/bin:/bin"
Environment="XDG_CONFIG_HOME={{ paperclip_home }}/.config"
Environment="XDG_CACHE_HOME={{ paperclip_home }}/.cache"

# Run Paperclip server
ExecStart={{ paperclip_bin }} run

# Restart policy
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ paperclip_home }}/.paperclip -{{ paperclip_home }}/.config/opencode -{{ paperclip_home }}/.cache/opencode -{{ paperclip_home }}/.local -{{ paperclip_home }}/.claude

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=paperclip

[Install]
WantedBy=multi-user.target
